# systemd/pi-log.service
# /etc/systemd/system/pi-log.service
[Unit]
Description=Pi-Log Beamrider Ingestion Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/pi-log

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="LOGEXP_URL=http://keep-0001.local:8000"
EnvironmentFile=-/etc/pi-log.env

# Command
ExecStart=/usr/bin/python3 -m pilog.main

# Restart Policy
Restart=always
RestartSec=2

# Security Hardening
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
